# deploy/systemd/highlight-manager.service
[Unit]
Description=Highlight Manager Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/highlight-manager
EnvironmentFile=/opt/highlight-manager/.env
ExecStart=/opt/highlight-manager/.venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
